# cbpp-configs/data/etc/profile.d/cbpp-gksu.sh
#!/bin/bash

alias gksu='pkexec --user root env DISPLAY=$DISPLAY XAUTHORITY=$XAUTHORITY'
alias gksudo='pkexec --user root env DISPLAY=$DISPLAY XAUTHORITY=$XAUTHORITY'
